# product-site/deploy/wenge-product-backup.timer
[Unit]
Description=Run Wenge product analytics backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=30m
Persistent=true
Unit=wenge-product-backup.service

[Install]
WantedBy=timers.target
